# docker/hexo/s6/ssh/finish
#!/bin/sh
sleep 5
echo "ssh server restarting..."
/usr/bin/sshd -D