# docker/hexo/s6/hexo/finish
#!/bin/sh
echo "hexo had ready" > /hexo.status